# lib/systemd.service
[Unit]
description = "Run the web app"

[Service]
Environment=/home/deployer/.bashrc
ExecStart=/env/bin/python3 manage.py runserver 0.0.0.0:$PORT
